# com.akihome.quant-evolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.akihome.quant-evolution</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/akihome/Documents/akimachome_antigravity/portfolio-bot/run_evolution.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/akihome/Documents/akimachome_antigravity/portfolio-bot</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 月曜 7:00 -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 水曜 7:00 -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 金曜 7:00 -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 土曜 9:00（週次サマリー通知） -->
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/quant-evolution-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/quant-evolution-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/tmp</string>
    </dict>

    <!-- Mac 起動時に初回実行（venv 自動構築含む） -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
